# cronrunpython.cron
*/30 * * * * /home/puja/Documents/sportskeeda-assignment/storeCSV.py